# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:trac]
command=/usr/local/bin/tracd -r --port 8000  --basic-auth="*,/root/test/passwd,Test Trac" /root/test
